# launchd/com.advisoryboard.memory-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.advisoryboard.memory-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/andraspolgar/Developer/advisory-board/scripts/memory_sync.py</string>
        <string>--persona</string>
        <string>hormozi</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/andraspolgar/Developer/advisory-board/memory/hormozi-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/andraspolgar/Developer/advisory-board/memory/hormozi-sync.error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/andraspolgar/Developer/advisory-board</string>
</dict>
</plist>
